# deploy/systemd/api-meusaldook.service.example
# Copiar para /etc/systemd/system/api-meusaldook.service e ajustar caminhos/variáveis.
#
#   sudo cp deploy/systemd/api-meusaldook.service.example /etc/systemd/system/api-meusaldook.service
#   sudo systemctl daemon-reload
#   sudo systemctl restart api-meusaldook
#
# Ajuste -Xmx conforme RAM da instância (ex.: t3.micro 1 GB → -Xmx384m; t3.small 2 GB → -Xmx768m).

[Unit]
Description=API Meu Saldo OK (Spring Boot)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/api-meusaldook

Environment=JAVA_HOME=/usr/lib/jvm/java-21-amazon-corretto
Environment=DATABASE_URL=jdbc:postgresql://127.0.0.1:5432/finance
Environment=DATABASE_USERNAME=finance
Environment=DATABASE_PASSWORD=finance
Environment=JWT_SECRET=altere-para-um-segredo-com-pelo-menos-32-bytes-utf8

# Metaspace: evita OutOfMemoryError: Metaspace em instâncias pequenas.
# ExitOnOutOfMemoryError: encerra o processo para o systemd reiniciar limpo.
Environment=JAVA_TOOL_OPTIONS=-Xms256m -Xmx512m -XX:MaxMetaspaceSize=256m -XX:+ExitOnOutOfMemoryError

ExecStart=/usr/bin/java -jar /home/ec2-user/api-meusaldook/target/api-meusaldook-0.0.1-SNAPSHOT.jar
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
